# dynamic-whitelist/systemd/gz-agent.service
[Unit]
Description=Dynamic whitelist Guangzhou nftables agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/dynamic-whitelist/gz-agent.env
WorkingDirectory=/opt/dynamic-whitelist/gz-agent
ExecStart=/usr/bin/python3 /opt/dynamic-whitelist/gz-agent/gz_agent.py
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/var/lib/dynamic-whitelist

[Install]
WantedBy=multi-user.target

