# dynamic-whitelist/systemd/hk-relay.service
[Unit]
Description=Dynamic whitelist HK relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/dynamic-whitelist/hk-relay.env
WorkingDirectory=/opt/dynamic-whitelist/hk-relay
ExecStart=/usr/bin/python3 /opt/dynamic-whitelist/hk-relay/hk_relay.py
Restart=always
RestartSec=3
User=dynamicwl
Group=dynamicwl
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target

